# i2pchat/blindbox/daemon/systemd/i2pchat-blindbox.service
[Unit]
Description=I2PChat BlindBox daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USER
WorkingDirectory=%h
EnvironmentFile=-%h/.config/i2pchat-blindbox/daemon.env
ExecStart=/usr/bin/python3 -m i2pchat.blindbox.daemon
Restart=always
RestartSec=2
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.i2pchat-blindbox %h/.config/i2pchat-blindbox
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
